# crates/cert/installfiles/systemd/ox_webservice.service
[Unit]
Description=ox Certificate Authority Server
Documentation=https://github.com/oxidizer/oxIDIZER
After=network-online.target
Wants=network-online.target

[Service]
Type=forking
User=ox_webservice
Group=ox_webservice
WorkingDirectory=/var/lib/ox_webservice
PIDFile=/var/run/ox_webservice/ox_webservice.pid

# Key passphrase - optionally set OX_CA_KEY_PASS in one of these files.
# Persistent:            echo 'OX_CA_KEY_PASS=your-passphrase' > /etc/ox_webservice/env && chmod 600 /etc/ox_webservice/env
# Runtime (no on-disk): sudo /usr/share/ox_webservice/unlock-ca.sh
EnvironmentFile=/etc/ox_webservice/env
EnvironmentFile=-/run/ox_webservice/env

# Verify the CA key is accessible before starting
ExecStartPre=+/usr/share/ox_webservice/check-ca-key.sh
ExecStart=/usr/share/ox_webservice/start_server.sh notice installed ca \
    /var/log/ox_webservice/ox_webservice.log \
    /var/run/ox_webservice/ox_webservice.pid
ExecStop=/usr/share/ox_webservice/stop_server.sh notice \
    /var/run/ox_webservice/ox_webservice.pid

Restart=on-failure
RestartSec=10

# Allow binding to ports 80 and 443 as non-root
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

# Restrict privileges
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/ox_webservice /var/log/ox_webservice /var/run/ox_webservice /etc/pki/ox_webservice/ca /etc/pki/ox_webservice/private/default /tmp
ReadOnlyPaths=/etc/ox_webservice /etc/pki/ox_webservice/private

[Install]
WantedBy=multi-user.target
